# baselines/deploy_asteroid/files/iperf_server.service
[Unit]
Description=iperf3 Server for Baselines Network Profiling
Documentation=https://iperf.fr/
After=network.target
# Never give up restarting
StartLimitIntervalSec=0

[Service]
Type=simple
# Kill any stale iperf3 processes hogging port 5201 before starting
ExecStartPre=/bin/sh -c 'fuser -k 5201/tcp 2>/dev/null; sleep 1; true'
# Persistent server mode — handles multiple sequential clients.
# The ansible liveness check force-restarts between serial profiling
# runs to clear any stuck/busy state.
# NOTE: --one-off was removed because it causes protocol mismatches
# on WAN links (test runs 15s instead of 5s, "Size of data read
# does not correspond to offered length" warning).
ExecStart=/usr/bin/iperf3 -s -p 5201
# Always restart if iperf3 crashes
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
